# speedtest/crontab
# basic 
*/31 * * * * bash /app/runSpeedtest.sh
