# speedtest/app/crontab
#
# basic 
*/31 * * * * bash /app/runSpeedtest.sh
@reboot sleep 15 && bash /app/runSpeedtest.sh
# 
#
